# SVCam/SVCam_Linux_x86/SVCamKit/others/svcamkit_u3v_startup.service
[Unit]
Description=U3V Startup parameter

[Service]
Type=simple
ExecStart=/etc/systemd/svcamkit_u3v_startup.sh

[Install]
WantedBy=multi-user.target
